# system/speaker-agent.service
[Unit]
Description=Bluetooth speaker agent

[Service]
ExecStart=python /home/pi/PiDPA/bin/speaker-agent.py

[Install]
WantedBy=default.target

